# infra/systemd/rustbrain-dev-health.timer
[Unit]
Description=Rustbrain dev-stack health check — every 5 minutes
Requires=rustbrain-dev-health.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=rustbrain-dev-health.service

[Install]
WantedBy=timers.target
