# infra/systemd/rustbrain-dev-health.service
[Unit]
Description=Rustbrain dev-stack health check
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot

# Adjust User= and WorkingDirectory= to match the actual clone path on mars.
User=jarnura
WorkingDirectory=/opt/rustbrain

ExecStart=/opt/rustbrain/scripts/dev-health-check.sh

# Exit 0 = healthy, exit 1 = failures detected (alert posted or cooldown active).
# Both are expected outcomes — treat as success so the timer keeps firing.
SuccessExitStatus=0 1

# Source tailscale port overrides into the health-check shell so
# CONTROL_API_HOST_PORT=18080, LOKI_HOST_PORT=13300, etc. resolve correctly.
# docker compose --env-file only exports vars to containers, not to this shell.
Environment="COMPOSE_ENV_FILE=/opt/rustbrain/compose/tailscale.env"

# Paperclip alerting credentials.
# Create /etc/rustbrain-dev-health.env on mars with:
#   PAPERCLIP_API_URL=https://<your-paperclip-host>
#   PAPERCLIP_API_KEY=<short-lived or long-lived token>
#   PAPERCLIP_COMPANY_ID=<uuid>
# The leading '-' means systemd ignores a missing file rather than failing.
EnvironmentFile=-/etc/rustbrain-dev-health.env

StandardOutput=journal
StandardError=journal
